# arch/airootfs/usr/lib/systemd/system/setup.service
[Unit]
Description=Install arch-server

[Service]
WorkingDirectory=/etc/setup
ExecStart=/bin/bash ./install.sh

[Install]
WantedBy=multi-user.target
